# deploy/systemd/ni-backend.service
[Unit]
Description=ni Companion AI Backend (FastAPI / uvicorn)
Documentation=https://github.com/lihaoalbert/ni
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 跟 ibiren 一致用 root 跑(国内 ECS 上避免 ni 用户 / .venv 权限问题)
# 真要降权:useradd ni && chown -R ni:ni /opt/ni 后改回 User=ni
User=root
Group=root
WorkingDirectory=/opt/ni/backend
EnvironmentFile=/opt/ni/backend/.env
Environment="PATH=/opt/ni/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/opt/ni/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --no-access-log

# 失败行为:非 0 退出就重启,5s 间隔,最多 5 次/分钟
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# 资源限制:防止 OOM
LimitNOFILE=65536
MemoryMax=1G

# 日志进 journald,可用 journalctl -u ni-backend 看
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ni-backend

# 安全加固(进程降权)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ni/backend/data /opt/ni/backend/.venv

[Install]
WantedBy=multi-user.target
